# templates/old/quadlet.container.j2
# Generated by Ansible. Do not edit.
[Unit]
Description={{ item.name }} container
Wants=network-online.target
After=network-online.target

[Container]
ContainerName={{ item.name }}
Image={{ item.image }}

{% for p in item.publish | default([]) %}
PublishPort={{ p }}
{% endfor %}

EnvironmentFile={{ env_dir }}/{{ item.name }}.env

{% for v in item.volumes | default([]) %}
Volume={{ v }}
{% endfor %}

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
